# systemd/threatcluster-ai-brief.service
[Unit]
Description=ThreatCluster AI Daily Brief Generator
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/threatcluster-cluster-service
Environment="PATH=/root/threatcluster-cluster-service/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/root/threatcluster-cluster-service"
ExecStart=/root/threatcluster-cluster-service/venv/bin/python -m src.ai_summary.daily_brief_scheduler
Restart=always
RestartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=threatcluster-ai-brief

# Resource limits
LimitNOFILE=65536
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target