# systemd/threatcluster-email.service
[Unit]
Description=ThreatCluster Daily Email Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/threatcluster-cluster-service
Environment="PATH=/root/threatcluster-cluster-service/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/root/threatcluster-cluster-service"

# Run the daily email scheduler
ExecStart=/root/threatcluster-cluster-service/venv/bin/python -m src.email_service.daily_email_scheduler

# Restart on failure
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=threatcluster-email

# Security settings - minimal for root user
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target